# deploy/examples/systemd/hwatu.example.service
[Unit]
Description=Hwatu API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hwatu
Group=hwatu
WorkingDirectory=/opt/hwatu/current
ExecStart=/usr/bin/dotnet /opt/hwatu/current/Hwatu.Server.dll
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5233
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
